# exam-scraper.service
[Unit]
Description=Exam Update Scraping System
After=network.target

[Service]
Type=simple
User=examscraper
Group=examscraper
WorkingDirectory=/opt/exam-scraper
Environment=PATH=/opt/exam-scraper/venv/bin
ExecStart=/opt/exam-scraper/venv/bin/python main.py --mode server --port 5000
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/exam-scraper/data /opt/exam-scraper/logs /opt/exam-scraper/backups

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
